# scibian-desktop/90-tg3.rules
ACTION=="add", SUBSYSTEM=="pci", ATTR{vendor}=="0x14e4", ATTR{device}=="0x1682", RUN+="/usr/sbin/scibian-fix-tg3"
